# etc/ex_servlet.service.in
[Unit]
Description=Start %%NAME%% servlet
Requires=network.target local-fs.target
After=fdaemons.service fdweb.service upsource.target

[Service]
Type=simple
RemainAfterExit=true
Environment=FD_LAUNCHER=%%NAME%%.service RUNUSER=@webuser@ RUNGROUP=@admin_group@ RUNDIR=@framerd_rundir@/servlets/ LOGDIR=@framerd_logdir@/servlets/
PassEnvironment=FD_LAUNCHER RUNUSER RUNGROUP
TimeoutStartSec=900
TimeoutStopSec=300
ExecStartPre=/bin/mkdir -p @framerd_rundir@/daemons/ @framerd_logdir@/daemons/ @framerd_logdir@/bugjar/
ExecStartPre=/bin/chown @fdaemon@ @framerd_rundir@/servlets/ @framerd_logdir@/servlets/ @framerd_logdir@/bugjar/
ExecStartPre=/bin/chgrp @admin_group@ @framerd_rundir@/servlets/ @framerd_logdir@/servlets/ @framerd_logdir@/bugjar/
ExecStartPre=/bin/chmod g+rwx @framerd_rundir@/servlets/ @framerd_logdir@/servlets/
ExecStartPre=/bin/chmod go+rwx @framerd_logdir@/bugjar/
ExecStart=@prefix@/bin/fdservlet @framerd_rundir@/servlets/%%NAME%.sock CONFIG=@config_dir@/config_files/%%NAME%%_servlet.cfg

[Install]
WantedBy=default.target
