# etc/fdaemons.service.in
[Unit]
Description=FramerD daemon starter
Requires=network.target local-fs.target
After=upsource.target

[Service]
Type=forking
RemainAfterExit=true
Environment=FD_LAUNCHER=framerd.service AS_USER=@webuser@ AS_GROUP=framerd
PassEnvironment=FD_LAUNCHER AS_USER AS_GROUP
ExecStartPre=@MKDIR@ -p @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHOWN@ @fdaemon@ @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHGRP@ @admin_group@ @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHMOD@ g+rwx @daemon_rundir@ @daemon_logdir@
ExecStart=@BINDIR@/fd_start_daemons @config_dir@/daemons
ExecStartPost=@TOUCHFILE@ @framerd_rundir@/daemons/_started
ExecStartPost=@DELETEFILE@ -f @framerd_rundir@/_sysv @framerd_rundir@/_upstart
ExecStartPost=@COPYFILE@ @framerd_rundir@/daemons/_started @framerd_rundir@/daemons/_systemd
ExecStop=@BINDIR@/fd_stop_jobs @daemon_rundir@
ExecStopPost=@DELETEFILE@ -f @framerd_rundir@/daemons/_started @framerd_rundir@/daemons/_systemd

[Install]
WantedBy=default.target
