# etc/fdweb.service.in
[Unit]
Description=FramerD servlet starter
Requires=network.target local-fs.target
After=fdaemons.service upsource.target

[Service]
Type=forking
RemainAfterExit=true
Environment=FD_LAUNCHER=fdweb.service RUNUSER=@webuser@ RUNGROUP=@webgroup@
PassEnvironment=FD_LAUNCHER RUNUSER RUNGROUP
TimeoutStartSec=300
TimeoutStopSec=300
ExecStartPre=@MKDIR@ -p @servlet_rundir@ @servlet_logdir@
ExecStartPre=@CHOWN@ @webuser@ @servlet_rundir@ @servlet_logdir@
ExecStartPre=@CHGRP@ @admin_group@ @servlet_rundir@ @servlet_logdir@
ExecStartPre=@CHMOD@ g+rwx @servlet_rundir@ @servlet_logdir@
ExecStart=@BINDIR@/fd_start_servlets @config_dir@/servlets
ExecStartPost=@TOUCHFILE@ @framerd_rundir@/servlets/_started
ExecStartPost=@DELETEFILE@ -f @framerd_rundir@/_sysv @framerd_rundir@/_upstart
ExecStartPost=@COPYFILE@ @framerd_rundir@/servlets/_started @framerd_rundir@/servlets/_systemd
ExecStop=@BINDIR@/fd_stop_jobs @servlet_rundir@
ExecStopPost=@DELETEFILE@ -f @framerd_rundir@/servlets/_started @framerd_rundir@/servlets/_systemd

[Install]
WantedBy=default.target
