# etc/FramerD.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>FramerD</string>
    <key>Description</key>
    <string>FramerD Server Manager</string>
    <key>ProgramArguments</key>
    <array>
      <string>@prefix@/bin/fdctl</string>
      <string>start</string>
    </array>
    <key>StandardOutPath</key>
    <string>@framerd_logdir@/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>@framerd_logdir@/launchd_error.log</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>AbandonProcessGroup</key>
    <true/>
  </dict>
</plist>
