# etc/framerd.conf.in
# FramerD upstart config, replaces fdmanager
start on (net-device-up and
          local-filesystems and
	  runlevel [345])
stop on runlevel [!2345]

author "Ken Haase"
description "FramerD servlets and server daemons"
expect none
console log

pre-start script
   mkdir -p @rundir@/fdserv
   mkdir -p @logdir@/fdserv
   mkdir -p @rundir@/framerd
   mkdir -p @logdir@/framerd
   chown @webuser@ @rundir@/fdserv @logdir@/fdserv
   chown @fdaemon@ @rundir@/framerd @logdir@/framerd
   chgrp @admin_group@ @rundir@/fdserv @logdir@/fdserv
   chgrp @admin_group@ @rundir@/framerd @logdir@/framerd
   chmod g+rwx @rundir@/fdserv  @logdir@/fdserv
   chmod g+rwx @rundir@/framerd @logdir@/framerd

   export FD_LOGLEVEL=7
   export COMMONCFG="STATEDIR=@rundir@/fdserv RESTART=1 PIDWAIT=no" 
   export SERVLETCFG="${COMMONCFG} RUNUSER=@webuser@" 
   export DAEMONCFG="${COMMONCFG} RUNGROUP=@admin_group@ RUNUSER=@fdaemon@" 
   for ctl in $(find -L @config_dir@/daemons/ -name "*.fdz");
     do if (readlink ${ctl}>/dev/null);
        then LOGFILE=@logdir@/framerd/`basename ${ctl} .fdz`.log \
	       fdserver `readlink ${ctl}` ${DAEMONCFG};
        else LOGFILE=@logdir@/framerd/`basename ${ctl} .fdz`.log \
	       fdserver ${ctl} ${DAEMONCFG}; 
        fi;
     done;
   for sock in $(find -L @config_dir@/servlets/ -name "*.sock");
     do if (readlink ${sock}>/dev/null);
        then touch @rundir@/fdserv/`basename ${sock}`.nospawn;
	     LOGFILE=@logdir@/fdserv/`basename ${sock} .sock`.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}` CONFIG=`readlink ${sock}` ${SERVLETCFG};
        else touch @rundir@/fdserv/`basename ${sock}`.nospawn;
	     LOGFILE=@logdir@/fdserv/`basename ${sock} .sock`.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}` CONFIG=${sock} ${SERVLETCFG};
        fi;
     done;
   for cfg in $(find -L @config_dir@/servlets/ -name "*.cfg");
     do if (readlink ${cfg}>/dev/null);
        then touch @rundir@/fdserv/`basename ${cfg} .cfg`.sock.nospawn;
	     LOGFILE=@logdir@/fdserv/`basename ${cfg} .cfg`.log \
	       fdservlet @rundir@/fdserv/`basename ${cfg} .cfg`.sock \
	       		 CONFIG=`readlink ${cfg}` ${SERVLETCFG};
        else touch @rundir@/fdserv/`basename ${sock} .cfg`.sock.nospawn;
	     LOGFILE=@logdir@/fdserv/`basename ${cfg} .cfg`.log \
	       fdservlet @rundir@/fdserv/`basename ${cfg} .cfg`.sock \
	       		 CONFIG=${sock} ${SERVLETCFG};
        fi;
     done;
   for sock in $(find -L @config_dir@/servlets/ -name "*@*");
     do if (readlink ${sock}>/dev/null);
        then LOGFILE=@logdir@/fdserv/${sock}.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}` CONFIG=`readlink ${sock}` ${SERVLETCFG};
        else LOGFILE=@logdir@/fdserv/${sock}.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}` CONFIG=${sock} ${SERVLETCFG};
	fi;
     done;
   for sock in $(find -L @config_dir@/servlets/ -name "*:*");
     do if (readlink ${sock}>/dev/null);
        then LOGFILE=@logdir@/fdserv/${sock}.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}` CONFIG=`readlink ${sock}` ${SERVLETCFG};
        else LOGFILE=@logdir@/fdserv/${sock}.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}` CONFIG=${sock} ${SERVLETCFG};
        fi;
     done;
end script

post-stop script
   for ppid in @rundir@/fdserv/*.ppid;
       do kill `cat $ppid` || true; 
          rm -f  @rundir@/fdserv/`basename $ppid .ppid`.nospawn
       done
   for ppid in  @rundir@/framerd/*.ppid;
       do kill `cat $ppid` || true; done
end script




