# etc/framerd-rc.d.in
#!/bin/sh
#
# Startup script for FramerD servers
# chkconfig:		- 80 80
# processname: 		fdmanager
# pidfile: 		/var/run/framerd.pid
### BEGIN INIT INFO
# Provides:		fdmanager
# Required-Start:	$local_fs $network $named $syslog
# Required-Stop:	$local_fs $network $named $syslog
# Default-Start:  	3 4 5
# Default-Stop:  	0 1 6
# Short-Description: 	Start configure FramerD servers
# Description: 		Start the FramerD servers described in etc/framerd/servers, including data and API servers
# X-Interactive:	false
### END INIT INFO

export FDSERVER=@prefix@/bin/fdserver
export FDSERVLET=@prefix@/bin/fdservlet
export FDAEMON_USER=@fdaemon@
export FDAEMON_GROUP=@admin_group@
export WEB_USER=@webuser@
CTLFILE=@config_dir@/servers
PIDFILE=/var/run/framerd.pid
STATFILE=/var/run/framerd.status

# See how we were called.
case "$1" in
  start)
	if test -f @daemon_rundir@/_started;
	then echo "FramerD daemons/servlets already started"; exit;
	else date > @daemon_rundir@/_started;
	fi
	for ctl in @config_dir@/daemons/*; do
	    if test ! -f ${ctl}; then echo > /dev/null;
	    elif (readlink ${ctl}>/dev/null);
	    then FD_RESTART=1 \
		 LOGFILE=@daemon_logdir@`basename ${ctl}`.log \
		 ${FDSERVER} `readlink ${ctl}` "PORT=`basename ${ctl}`" \
		             STATEDIR=@daemon_rundir@;
	    else FD_RESTART=1 \
		 LOGFILE=@daemon_logdir@`basename ${ctl}`.log \
		 ${FDSERVER} ${ctl} "PORT=`basename ${ctl}`" \
		             STATEDIR=@daemon_rundir@;
	    fi;
	done;
	for ctl in @config_dir@/servlets/*@* @config_dir@/servlets/*:*; do
	    if test ! -f ${ctl}; then echo > /dev/null;
	    elif (readlink ${ctl}>/dev/null);
	    then FD_RESTART=1 \
		 LOGFILE=@servlet_logdir@/`basename ${ctl}`.log \
		 ${FDSERVLET} `basename ${ctl}` "CONFIG=`readlink ${ctl}`" \
		             STATEDIR=@servlet_rundir@;
	    else FD_RESTART=1 \
		 LOGFILE=@servlet_logdir@/`basename ${ctl}`.log \
		 ${FDSERVLET} `basename ${ctl}` "CONFIG=${ctl}" \
		             STATEDIR=@servlet_rundir@;
	    fi;
	done;
	for ctl in @config_dir@/servlets/*.sock; do
	    if test ! -f ${ctl}; then echo > /dev/null;
	    elif (readlink ${ctl}>/dev/null);
	    then FD_RESTART=1 \
		 LOGFILE=@servlet_logdir@/`basename ${ctl} .sock`.log \
		 ${FDSERVLET} `basename ${ctl}` "CONFIG=`readlink ${ctl}`" \
		             STATEDIR=@servlet_rundir@;
	    else FD_RESTART=1 \
		 LOGFILE=@servlet_logdir@/`basename ${ctl} .sock`.log \
		 ${FDSERVLET} `basename ${ctl}` "CONFIG=${ctl}" \
		             STATEDIR=@servlet_rundir@;
	    fi;
	done;
	;;
  stop)
	for ctl in @config_dir@/daemons/*; do
	    ppid_file=@daemon_rundir@/`basename ${ctl}`.ppid
	    if test ! -f ${ppid_file}; then echo > /dev/null;
	    elif test @daemon_rundir@/_started -nt ${ppid_file};
	    then rm ${ppid_file};
	    else kill `cat ${ppid_file}`;
	    fi;
	done;
	for ctl in @config_dir@/servlets/*@* @config_dir@/servlets/*:*; do
	    ppid_file=@servlet_rundir@/`basename ${ctl}`.ppid
	    if test ! -f ${ppid_file}; then echo > /dev/null;
	    elif test @daemon_rundir@/_started -nt ${ppid_file};
	    then rm ${ppid_file};
	    else kill `cat ${ppid_file}`;
	    fi;
	done;
	for ctl in @config_dir@/servlets/*.sock; do
	    ppid_file=@servlet_rundir@/`basename ${ctl} .sock`.ppid
	    if test ! -f ${ppid_file}; then echo > /dev/null;
	    elif test @daemon_rundir@/_started -nt ${ppid_file};
	    then rm ${ppid_file};
	    else kill `cat ${ppid_file}`;
	    fi;
	done;
	rm  @daemon_rundir@/_started;
	;;
  status)
	if test -f @daemon_rundir@/_started;
	    echo -n "Daemons: "; ls @daemon_rundir@/*.ppid;
	    echo -n "Servlets: "; ls @servlet_rundir@/*.ppid;
	fi
	;;
  restart|force-reload)
	$0 stop
	echo "Sleeping for 5 seconds before restart"
	sleep 5
	$0 start
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
esac

exit 0
