# etc/framerd.conf.in
# FramerD upstart config, replaces fdmanager
start on (net-device-up and
          local-filesystems and
	  runlevel [345])
stop on runlevel [!2345]

author "Ken Haase"
description "FramerD servlets and server daemons"
expect none
console log

pre-start script
   mkdir -p @rundir@/fdserv
   mkdir -p @logdir@/fdserv
   mkdir -p @rundir@/framerd
   mkdir -p @logdir@/framerd
   chown @webuser@ @rundir@/fdserv @logdir@/fdserv
   chown @fdaemon@ @rundir@/framerd @logdir@/framerd
   chgrp @admin_group@ @rundir@/fdserv @logdir@/fdserv
   chgrp @admin_group@ @rundir@/framerd @logdir@/framerd
   chmod g+rwx @rundir@/fdserv  @logdir@/fdserv
   chmod g+rwx @rundir@/framerd @logdir@/framerd

   for ctl in $(find @config_dir@/daemons/ -name "*.fdz" \
       	             -mindepth 1 -maxdepth 3);
     do if (readlink ${ctl}>/dev/null);
        then LOGFILE=@logdir@/framerd/`basename ${ctl} .fdz`.log	\
	       fdserver `readlink ${ctl}` STATEDIR=@rundir@/framerd	\
	          GROUP=@admin_group@ USER=@fdaemon@ RESTART=1;
        else LOGFILE=@logdir@/framerd/`basename ${ctl} .fdz`.log	\
               fdserver ${ctl} STATEDIR=@rundir@/framerd		\
	       	  GROUP=@admin_group@ USER=@fdaemon@ RESTART=1; 
        fi;
     done;
   for sock in $(find @config_dir@/servlets/ -name "*.sock" 		\
                      -mindepth 1 -maxdepth 3);
     do if (readlink ${sock}>/dev/null);
        then touch @rundir@/fdserv/`basename ${sock}`.nospawn;
	     LOGFILE=@logdir@/fdserv/`basename ${sock} .sock`.log	\
               fdservlet @rundir@/fdserv/`basename ${sock}`		\
	         CONFIG=`readlink ${sock}` STATEDIR=@rundir@/fdserv	\ 
		 GROUP=@admin_group@ USER=@webuser@			\
		 RESTART=1 PIDWAIT=no;
        else touch @rundir@/fdserv/`basename ${sock}`.nospawn;
	     LOGFILE=@logdir@/fdserv/`basename ${sock} .sock`.log	\
 	       fdservlet @rundir@/fdserver/`basename ${sock}`		\ 
 	         CONFIG=${sock} STATEDIR=@rundir@/fdserv 		\
	         GROUP=@admin_group@ USER=@webuser@ 			\
		 RESTART=1;
        fi;
     done;
   for sock in $(find @config_dir@/servlets/ -name "*@*"
                      -mindepth 1 -maxdepth 3);
     do if (readlink ${sock}>/dev/null);
        then LOGFILE=@logdir@/fdserv/${sock}.log \
	       fdservlet @rundir@/fdserv/`basename ${sock}`		\
	         CONFIG=`readlink ${sock}` STATEDIR=@rundir@/fdserv 	\
	         GROUP=@admin_group@ USER=@webuser@          		\
		 PIDWAIT=no RESTART=1;
        else LOGFILE=@logdir@/fdserv/${sock}.log		\
	       fdservlet @rundir@/fdserv/`basename ${sock}`	\
	         CONFIG=${sock} STATEDIR=@rundir@/fdserv	\
		  GROUP=@admin_group@ USER=@webuser@		\
		 PIDWAIT=no RESTART=1;
	fi;
     done;
   for sock in $(find @config_dir@/servlets/ -name "*:*" \
                      -mindepth 1 -maxdepth 3);
     do if (readlink ${sock}>/dev/null);
        then LOGFILE=@logdir@/fdserv/${sock}.log 			\
	       fdservlet @rundir@/fdserv/`basename ${sock}`		\
	         CONFIG=`readlink ${sock}` STATEDIR=@rundir@/fdserv 	\
		 GROUP=@admin_group@ USER=@webuser@			\
		 PIDWAIT=no RESTART=1;
        else LOGFILE=@logdir@/fdserv/${sock}.log 		\
	       fdservlet @rundir@/fdserv/`basename ${sock}`	\
	         CONFIG=${sock} STATEDIR=@rundir@/fdserv	\
		 GROUP=@admin_group@ USER=@webuser@		\
		 PIDWAIT=no RESTART=1;
        fi;
     done;
end script

post-stop script
   for ppid in @rundir@/fdserv/*.ppid;
       do kill `cat $ppid` || true; 
          rm -f  @rundir@/fdserv/`basename $ppid .ppid`.nospawn
       done
   for ppid in  @rundir@/framerd/*.ppid;
       do kill `cat $ppid` || true; done
end script




