# etc/framerd.conf.in
# FramerD upstart config, replaces fdmanager
start on (net-device-up and
          local-filesystems and
	  runlevel [345])
stop on runlevel [!2345]

author "Ken Haase"
description "FramerD servlets and server daemons"
expect none
console log

pre-start script
   mkdir -p @servlet_rundir@
   mkdir -p @logdir@/fdserv
   mkdir -p @daemon_rundir@
   mkdir -p @logdir@/framerd
   chown @webuser@ @servlet_rundir@ @logdir@/fdserv
   chown @fdaemon@ @daemon_rundir@ @logdir@/framerd
   chgrp @admin_group@ @servlet_rundir@ @logdir@/fdserv
   chgrp @admin_group@ @daemon_rundir@ @logdir@/framerd
   chmod g+rwx @servlet_rundir@  @logdir@/fdserv
   chmod g+rwx @daemon_rundir@ @logdir@/framerd

   export FD_LAUNCHER=framerd.conf
   fdstart @config_dir@/daemons/
   fdstart @config_dir@/servlets/
end script

post-stop script
   fdstop @servlet_rundir@ @daemon_rundir@
end script
