# etc/framerd.service.in
[Unit]
Description=FramerD service(s) starter
Requires=network.target
After=upsource.service cloud-final.service

[Service]
Type=simple
RemainAfterExit=true
Environment=FD_LAUNCHER=framerd.service
PassEnvironment=FD_LAUNCHER
ExecStartPre=@MKDIR@ -p @daemon_rundir@ @daemon_logdir@ @servlet_rundir@ @servlet_logdir@
ExecStartPre=@CHOWN@ @webuser@ @servlet_rundir@ @servlet_logdir@
ExecStartPre=@CHOWN@ @fdaemon@ @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHGRP@ @admin_group@ @servlet_rundir@ @servlet_logdir@
ExecStartPre=@CHGRP@ @admin_group@ @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHMOD@ g+rwx @servlet_rundir@  @servlet_logdir@
ExecStartPre=@CHMOD@ g+rwx @daemon_rundir@ @daemon_logdir@
ExecStart=@BINDIR@/fdstart @config_dir@/daemons @config_dir@/servlets
ExecStartPost=@TOUCHFILE@ @framerd_rundir@/_started
ExecStartPost=@DELETEFILE@ -f @framerd_rundir@/_sysv @framerd_rundir@/_upstart
ExecStartPost=@COPYFILE@ @framerd_rundir@/_started @framerd_rundir@/_systemd
ExecStop=@BINDIR@/fdstop @daemon_rundir@ @servlet_rundir@
ExecStopPost=@DELETEFILE@ -f @framerd_rundir@/_started @framerd_rundir@/_systemd

[Install]
WantedBy=multi-user.target
