# etc/framerd.target.in
[Unit]
Description=FramerD service(s) starter
Requires=fdaemons.service fdservlets.service
After=upsource.target

[Install]
WantedBy=default.target
