# etc/framerd.target
[Unit]
Description=FramerD service(s) starter
Requires=fdaemons.service fdweb.service
After=upsource.target

[Install]
WantedBy=default.target
