# etc/framerd.conf.in
# FramerD upstart config, replaces fdmanager
start on (net-device-up and
          local-filesystems and
	  runlevel [345])
stop on runlevel [!2345]

author "Ken Haase"
description "FramerD servlets and server daemons"
expect none
console log

pre-start script
   mkdir -p @rundir@/fdserv
   mkdir -p @logdir@/fdserv
   mkdir -p @rundir@/framerd
   mkdir -p @logdir@/framerd
   chown @webuser@ @rundir@/fdserv @logdir@/fdserv
   chown @fdaemon@ @rundir@/framerd @logdir@/framerd
   chgrp @admin_group@ @rundir@/fdserv @logdir@/fdserv
   chgrp @admin_group@ @rundir@/framerd @logdir@/framerd
   chmod g+rwx @rundir@/fdserv  @logdir@/fdserv
   chmod g+rwx @rundir@/framerd @logdir@/framerd

   export FD_LOGLEVEL=7
   fdstart @config_dir@/daemons/
   fdstart @config_dir@/servlets/
end script

post-stop script
   fdstop @rundir@/fdserv;
   fdstop @rundir@/framerd
end script




