# etc/kno-rc.d.in
#!/bin/sh
#
# Startup script for Kno servers
# chkconfig:		- 80 80
# processname: 		kno
# pidfile: 		/var/run/kno.pid
### BEGIN INIT INFO
# Provides:		kno
# Required-Start:	$local_fs $network $named $syslog
# Required-Stop:	$local_fs $network $named $syslog
# Default-Start:  	3 4 5
# Default-Stop:  	0 1 6
# Short-Description: 	Start configure Kno servers
# Description: 		Start the Kno servers described in etc/kno/servers, including data and API servers
# X-Interactive:	false
### END INIT INFO

export KNOD=@prefix@/bin/knod
export KNOWEB=@prefix@/bin/knoweb
export KNODAEMON_USER=@knoaemon@
export KNODAEMON_GROUP=@admin_group@
export WEB_USER=@webuser@
CTLFILE=@config_dir@/servers
PIDFILE=/var/run/kno.pid
STATFILE=/var/run/kno.status

# See how we were called.
case "$1" in
  start)
	if test -f @kno_rundir@/_started;
	then echo "Kno daemons/servlets already started"; exit;
	else date > @kno_rundir@/_started; 
	     rm -f @kno_rundir@/_upstart @kno_rundir@/_system;
	     cp @kno_rundir@/_started @kno_rundir@/_sysv;
	fi
	knostartall;
	;;
  stop)
	knostopall;
	rm -f @daemon_rundir@/_started @daemon_rundir@/_sysv;
	;;
  status)
	if test -f @kno_rundir@/_started;
	    echo "Started: " `cat @kno_rundir@/_started`
	    echo -n "Daemons: "; ls @daemon_rundir@/*.ppid;
	    echo -n "Servlets: "; ls @servlet_rundir@/*.ppid;
	fi
	;;
  restart|force-reload)
	$0 stop
	echo "Sleeping for 5 seconds before restart"
	sleep 5
	$0 start
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status}"
	exit 1
esac

exit 0
