# etc/upstart/kno.conf.in
# Kno upstart config, replaces fdmanager
start on (net-device-up and
          local-filesystems and
	  runlevel [345])
stop on runlevel [!2345]

author "Ken Haase"
description "Kno servlets and server daemons"
expect none
console log

pre-start script
   mkdir -p @servlet_rundir@
   mkdir -p @daemon_logdir@
   mkdir -p @servlet_logdir@
   mkdir -p @daemon_logdir@
   chown @webuser@ @servlet_rundir@ @servlet_logdir@
   chown @knouser@ @daemon_rundir@ @daemon_logdir@
   chgrp @admin_group@ @servlet_rundir@
   chgrp @admin_group@ @daemon_rundir@
   chmod g+rwx @servlet_rundir@  @servlet_logdir@
   chmod g+rwx @daemon_rundir@ @daemon_logdir@

   export KNO_LAUNCHER=kno.conf
   fdstart @config_dir@/daemons/;
   fdstart @config_dir@/servlets/;
   date > @kno_rundir@/_started; 
   rm -f @kno_rundir@/_sysv @kno_rundir@/_systemd;
   cp @kno_rundir@/_started @kno_rundir@/_upstart;
end script

post-stop script
   fdstop @servlet_rundir@ @daemon_rundir@;
   rm -f @daemon_rundir@/_started @daemon_rundir@/_sysv;
end script
