# dist/debian/systemd-kno.conffiles
/lib/systemd/system/kno.target
/lib/systemd/system/kno-daemons.service
/lib/systemd/system/kno-servlets.service
