# etc/systemd/ex_servlet.service.in
[Unit]
Description=Start %%NAME%% servlet
Requires=network.target local-fs.target
After=kno-daemons.service kno-servlets.service upsource.target

[Service]
Type=simple
RemainAfterExit=true
Environment=KNO_LAUNCHER=%%NAME%%.service RUNUSER=@webuser@ RUNGROUP=@admin_group@ RUNDIR=@kno_rundir@/servlets/ LOGDIR=@kno_logdir@/servlets/
PassEnvironment=KNO_LAUNCHER RUNUSER RUNGROUP
TimeoutStartSec=900
TimeoutStopSec=300
ExecStartPre=/bin/mkdir -p @kno_rundir@/daemons/ @kno_logdir@/daemons/ @kno_logdir@/bugjar/
ExecStartPre=/bin/chown @knodaemon@ @kno_rundir@/servlets/ @kno_logdir@/servlets/ @kno_logdir@/bugjar/
ExecStartPre=/bin/chgrp @admin_group@ @kno_rundir@/servlets/ @kno_logdir@/servlets/ @kno_logdir@/bugjar/
ExecStartPre=/bin/chmod g+rwx @kno_rundir@/servlets/ @kno_logdir@/servlets/
ExecStartPre=/bin/chmod go+rwx @kno_logdir@/bugjar/
ExecStart=@prefix@/bin/knocgi @kno_rundir@/servlets/%%NAME%.sock CONFIG=@config_dir@/config_files/%%NAME%%_servlet.cfg

[Install]
WantedBy=default.target
