# etc/systemd/kno-daemons.service.in
[Unit]
Description=Kno daemon starter
Requires=network.target local-fs.target
After=upsource.target

[Service]
Type=forking
RemainAfterExit=true
Environment=KNO_LAUNCHER=kno.service RUNUSER=@knouser@ RUNGROUP=@admin_group@
PassEnvironment=KNO_LAUNCHER RUNUSER=@knouser@ RUNGROUP=@admin_group@
TimeoutStartSec=300
TimeoutStopSec=300
ExecStartPre=@MKDIR@ -p @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHOWN@ @knouser@ @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHGRP@ @admin_group@ @daemon_rundir@ @daemon_logdir@
ExecStartPre=@CHMOD@ g+rwx @daemon_rundir@ @daemon_logdir@
ExecStart=@BINDIR@/kno_start_daemons @config_dir@/daemons
ExecStartPost=@TOUCHFILE@ @kno_rundir@/daemons/_started
ExecStartPost=@DELETEFILE@ -f @kno_rundir@/_sysv @kno_rundir@/_upstart
ExecStartPost=@COPYFILE@ @kno_rundir@/daemons/_started @kno_rundir@/daemons/_systemd
ExecStop=@BINDIR@/kno_stop_jobs @daemon_rundir@
ExecStopPost=@DELETEFILE@ -f @kno_rundir@/daemons/_started @kno_rundir@/daemons/_systemd

[Install]
WantedBy=default.target
