# etc/systemd/kno.target
[Unit]
Description=KNO service(s) starter
Requires=kno-daemons.service kno-servlets.service
After=upsource.target

[Install]
WantedBy=default.target
